# .docker/etc/supervisor.d/caddy.ini
[group:webserver]
programs=caddy
priority=20

[program:caddy]
process_name=%(program_name)s_%(process_num)02d
command=caddy -conf %(ENV_DOCKER_ETC_DIR)s/caddy/Caddyfile
startsecs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0